# deployment/systemd/dart-backend.service
[Unit]
Description=Dart Tournament Backend API
After=network.target postgresql.service redis-server.service
Wants=postgresql.service redis-server.service

[Service]
Type=simple
User=darts-admin
Group=darts-admin
WorkingDirectory=/home/darts-admin/DartTournament
Environment="PYTHONPATH=/home/darts-admin/DartTournament"
ExecStart=/home/darts-admin/DartTournament/backend/venv/bin/uvicorn backend.main:app --host 0.0.0.0 --port 8000
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
