# deployment/systemd/dart-display.service
[Unit]
Description=Dart Tournament Display Terminal
After=network.target dart-backend.service
Wants=dart-backend.service

[Service]
Type=simple
User=darts-admin
Group=darts-admin
WorkingDirectory=/home/darts-admin/DartTournament/display-terminal
Environment="NODE_ENV=production"
ExecStart=/home/darts-admin/DartTournament/node_modules/.bin/next start -p 3002
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
